# deploy/lexgraph-refresh.timer
[Unit]
Description=Daily Lexgraph refresh at 04:37 UTC (off-peak)

[Timer]
OnCalendar=*-*-* 04:37:00 UTC
Persistent=true
RandomizedDelaySec=8m

[Install]
WantedBy=timers.target
